# ops/systemd/solana-backup-sqlite.service
[Unit]
Description=Solana Trader — SQLite Hot Backup
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/root/solana_trader/ops/backup_sqlite.sh 15min
StandardOutput=journal
StandardError=journal
